# cloud-init/etc/profile.d/ceph_default.sh
#!/bin/sh
export CEPH_ARGS="--id smartaqnet"
